# deploy/systemd/backup.timer
[Unit]
Description=Backup diario del panel

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
