# deploy/systemd/monitor.timer
[Unit]
Description=Chequeo de salud del panel cada 2 minutos

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
